# crontab
0 * * * * python3 /app/crawl.py >> /var/log/cron.log 2>&1
